# booterfilter/booterfilter.cron
0 0 1 * *	nemead	/usr/bin/nemea/booterfilter-update.sh
